# supervisor/kinaya-horizon.conf
; ─────────────────────────────────────────────────────────────────────────────
; Supervisor configuration for Kinaya POS — Laravel Horizon
;
; INSTALLATION (Ubuntu / Debian):
;   sudo apt-get install supervisor
;   sudo cp supervisor/kinaya-horizon.conf /etc/supervisor/conf.d/
;   sudo supervisorctl reread
;   sudo supervisorctl update
;   sudo supervisorctl start kinaya-horizon:*
;
; Horizon manages ALL queue workers (default + webhooks) internally.
; You only need ONE supervisor process pointing at `php artisan horizon`.
; ─────────────────────────────────────────────────────────────────────────────

[program:kinaya-horizon]
; ── Adjust paths to match your deployment ────────────────────────────────────
process_name=%(program_name)s
command=php /var/www/kinaya/backend/artisan horizon
autostart=true
autorestart=true
; Restart horizon gracefully when the process exits (e.g., after deploy)
; Send SIGTERM to let Horizon finish in-flight jobs before stopping
stopsignal=SIGTERM
stopwaitsecs=3600

; ── User running the process (should own the project files) ──────────────────
user=www-data
; On shared hosting or VPS where the deploy user differs, change to that user.

; ── Output ───────────────────────────────────────────────────────────────────
redirect_stderr=true
stdout_logfile=/var/www/kinaya/backend/storage/logs/horizon-supervisor.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile=/var/www/kinaya/backend/storage/logs/horizon-supervisor-err.log

; ── Environment variables injected into the process ──────────────────────────
; Duplicate only keys that differ from the .env file, or leave empty.
environment=APP_ENV="production"
